# multi-monitor-wallpaper.service
[Unit]
Description=Multi-Monitor Wallpaper Cycler
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/marcus_aurelius/multi-monitor-wallpaper.py
Restart=always
RestartSec=10
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/marcus_aurelius/.Xauthority"

[Install]
WantedBy=default.target