# multi-monitor-wallpaper.service
[Unit]
Description=Multi-Monitor Wallpaper Cycler
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 $HOME/multi-monitor-wallpaper.py
Restart=always
RestartSec=10
Environment="DISPLAY=:0"
Environment="XAUTHORITY=$HOME/.Xauthority"

[Install]
WantedBy=default.target