# multi-monitor-wallpaper.service
[Unit]
Description=Multi-Monitor Wallpaper Cycler
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/multi-monitor-wallpaper.py
Restart=always
RestartSec=10
Environment=DISPLAY=:0
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=XDG_RUNTIME_DIR=%h/.config

[Install]
WantedBy=default.target